# scripts/com.jyukipann.tts-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jyukipann.tts-api</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/path/to/tts-api/.venv/bin/uvicorn</string>
        <string>tts_api.app:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8010</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/path/to/tts-api</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HF_HOME</key>
        <string>/Users/YOUR_USERNAME/path/to/tts-api/.cache/huggingface</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/path/to/tts-api/server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/path/to/tts-api/server.log</string>
</dict>
</plist>
